# tests/unit/openwrt/golden/ath79/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# TODO make this a true service
/root/bin/wifi-details.sh >/dev/null 2>/dev/null </dev/null &

# By default anything install in opkg that includes a service will start out
# of the box. Easiest way to disable services that we depend
# on but dont want to start
/etc/init.d/uhttpd stop; /etc/init.d/uhttpd disabled

# apinger initial template population
# this needs to remain in the case where DHCP isnt working initially
# there is also a race condition with the wifi interfaces coming online
# and requires a 30 sec grace period before checking via apinger. If not
# the result is wifi is up but the wifi tool detects it as down
(sleep 30 && /root/bin/apinger-pop.sh "8.8.8.8") &> /dev/null &

exit 0
