# tests/unit/openwrt/golden/ar71xx/etc/crontabs/root
0 * * * * /usr/sbin/logrotate -f /etc/logrotate.conf -s /tmp/logrotate.status
