# tests/unit/openwrt/golden/ar71xx/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# TODO make this a true service
/root/bin/wifi-details.sh >/dev/null 2>/dev/null </dev/null &

exit 0
